# roles/prometheus/templates/prometheus.service.j2
# =============================================================================
# templates/prometheus.service.j2 - SystemD service template
# =============================================================================

# {{ ansible_managed }}
# Prometheus SystemD service file
# Generated by Ansible prometheus role

[Unit]
Description=Prometheus Time Series Collection and Processing Server
Documentation=https://prometheus.io/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User={{ prometheus_system_user }}
Group={{ prometheus_system_group }}
ExecStart={{ prometheus_binary_install_dir }}/prometheus \
  --config.file={{ prometheus_config_file }} \
  --storage.tsdb.path={{ prometheus_storage_tsdb_path }} \
  --storage.tsdb.retention.time={{ prometheus_storage_retention }} \
  --storage.tsdb.retention.size={{ prometheus_storage_retention_size }} \
{% if prometheus_storage_tsdb_wal_compression %}
  --storage.tsdb.wal-compression \
{% endif %}
  --web.listen-address={{ prometheus_web_listen_address }} \
  --web.telemetry-path={{ prometheus_web_telemetry_path }} \
{% if prometheus_web_external_url %}
  --web.external-url={{ prometheus_web_external_url }} \
{% endif %}
{% if prometheus_web_route_prefix != '/' %}
  --web.route-prefix={{ prometheus_web_route_prefix }} \
{% endif %}
{% if prometheus_web_config | length > 0 %}
  --web.config.file={{ prometheus_web_config_file }} \
{% endif %}
  --web.console.templates={{ prometheus_config_dir }}/consoles \
  --web.console.libraries={{ prometheus_config_dir }}/console_libraries \
  --web.enable-lifecycle \
  --log.level=info \
  --log.format=logfmt
{% for flag, value in prometheus_config_flags_extra.items() %}
  --{{ flag }}={{ value }} \
{% endfor %}

ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=20s
SendSIGKILL=no
Restart={{ prometheus_systemd_restart_policy }}
RestartSec={{ prometheus_systemd_restart_sec }}

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths={{ prometheus_data_dir }} {{ prometheus_log_dir }}
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes

# Resource limits
LimitNOFILE={{ prometheus_limit_nofile }}
LimitMEMLOCK={{ prometheus_limit_memlock }}

# Working directory
WorkingDirectory={{ prometheus_data_dir }}

# Runtime directory
RuntimeDirectory=prometheus
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target